# linux/70-stm32f0-pmic-beeper.rules
ATTRS{name}=="PC Speaker", ENV{DEVNAME}!="", TAG+="uaccess", RUN+="/usr/bin/stm32f0-beeper-alias"
